# services/waze-sidecar/deploy/waze-sidecar.service
[Unit]
Description=Waze Playwright sidecar
After=network.target xvfb.service
Requires=xvfb.service

[Service]
Type=simple
User=waze-sidecar
Group=waze-sidecar
WorkingDirectory=/opt/waze-sidecar
EnvironmentFile=/opt/waze-sidecar/.env
Environment=HOME=/opt/waze-sidecar
Environment=DISPLAY=:99
Environment=PLAYWRIGHT_HEADLESS=0
Environment=COOKIE_FILE=/opt/waze-sidecar/data/waze_session.json
Environment=PLAYWRIGHT_BROWSERS_PATH=/opt/waze-sidecar/.playwright-browsers
ExecStart=/opt/waze-sidecar/.venv/bin/python -m uvicorn app:app --host 127.0.0.1 --port 8001
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=waze-sidecar

[Install]
WantedBy=multi-user.target
